# ansible/files/scripts/rc.local
#!/bin/bash

# Preventing instant resume after suspend
echo XHC1 > /proc/acpi/wakeup
echo LID0 > /proc/acpi/wakeup
