# ansible/files/scripts/rc.local
#!/bin/sh -e

# disabling wakeup
echo "PXSX" > /proc/acpi/wakeup
echo "GLAN" > /proc/acpi/wakeup
echo "XHC" > /proc/acpi/wakeup

exit 0